# supervisor.conf
[program:patio]
command=/home/ortiz/www/patio/start.sh
user=ortiz
stdout_logfile=/home/ortiz/www/patio/logs/supervisor.log
redirect_stderr=true
autostart=true
autorestart=true
redirect_stderr=true
